# capture/capture.cron
* * * * * bash /home/mspitz/src/eckleburg/capture/capture_single.sh
